# Install/system-files/common/etc/systemd/system/check-bat.service
[Unit]
Description=Check battery level and suspend if low

[Service]
Type=simple
ExecStart=/usr/local/bin/check-bat
