# Install/system-files/common/etc/systemd/system/check-bat.timer
[Unit]
Description=Run check-battery service every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
